# supervisor_mac.ini
[program:django-reference]
command=sphinx-autobuild
    --color
    --host 0.0.0.0
    -j auto
    --port 18001
    . _build/html/
    --re-ignore "\.mypy_cache"
    --re-ignore "\.git"
    --re-ignore "\.*\.swp"
    --re-ignore "\.*\.log"

directory=/Users/wangx/github/django/django-reference
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=HOME="/Users/wangx",USER="wangx",PATH="/Users/wangx/venv/bin/:/Users/wangx/.local/bin:/Users/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin",RAMWIN="DJANGO_REFERENCE"
stdout_logfile=/Users/wangx/github/django/django-reference/stdout_%(process_num)d.log
stderr_logfile=/Users/wangx/github/django/django-reference/stderr_%(process_num)d.log
