# supervisor.ini
[program:django-reference]
command=sh run_server.sh
directory=/home/wangx/github/django/django-reference
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=20
stderr_logfile_maxbytes=4MB
stderr_logfile_backups=20
autostart=true
autorestart=true
startretries=3
startsecs=20
redirect_stderr=false
user=wangx
environment=PATH="/home/wangx/venv/bin/:/home/wangx/.local/bin:/home/wangx/node/bin/:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
stdout_logfile=/home/wangx/github/django/django-reference/stdout.log
stderr_logfile=/home/wangx/github/django/django-reference/stderr.log
